# backups/crontab
0 2 * * * pg_dump -h db -U postgres securedb > /backups/securedb_$(date +\%F).sql
